# sv/snooze-daily/finish
#!/bin/sh
exec touch /var/cache/snooze/daily
